# aw-notion.service
[Unit]
Description=Sync ActivityWatch focus blocks to Notion Time Log
After=network-online.target

[Service]
Type=oneshot
ExecStart=BIN_PATH sync
