# com.aw-watcher-window.keepalive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aw-watcher-window.keepalive</string>
    <key>ProgramArguments</key>
    <array>
        <string>BIN_PATH</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>LOG_DIR/window-keepalive.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_DIR/window-keepalive.log</string>
</dict>
</plist>
